# brainstorm_refresh.timer
[Unit]
Description=Timer

[Timer]
OnUnitActiveSec=120s
OnBootSec=120s

[Install]
WantedBy=timers.target